# deployment/md_web/md_web.service
[Unit]
Description=md-web (Obsidian-style vault)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/cow/repos/md-web
ExecStart=/usr/bin/env PORT=3009 NODE_ENV=production /usr/bin/node server.js
EnvironmentFile=-/home/cow/repos/md-web/.env
Restart=always

[Install]
WantedBy=multi-user.target
